# conf/supervisord/phantom.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[group:phantom]
programs=node,server

;start phantom-node
[program:node]
command=node index.js
autostart=true
autorestart=true
startretries=3
user=InstantForge
loglevel=info
directory=/madeira/source/html5/phantom
stdout_logfile=/madeira/log/supervisord/phantom-node.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true

;start phantom-server
[program:server]
command=phantomjs phantom_server.js
autostart=true
autorestart=true
startretries=3
user=InstantForge
loglevel=info
directory=/madeira/source/html5/phantom
stdout_logfile=/madeira/log/supervisord/phantom-server.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
