# conf/supervisord/requestworker.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[eventlistener:RequestWorker]
command=python -u /usr/bin/rubber -p requestworker* -n 8
events=SUPERVISOR_STATE_CHANGE_RUNNING
stdout_events_enabled=false
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs
stdout_logfile=/visualops/log/supervisord/reqestworker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true

;start requestworker
[program:requestworker1]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker1.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker2]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker2.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker3]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker3.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker4]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker4.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker5]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker5.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker6]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
redirect_stderr=true
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker6.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker7]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker7.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/

[program:requestworker8]
command=python26 /visualops/source/api/Source/INiT/Instant/Forge/RequestWorker/RequestWorker.py
autorestart=true
startretries=3
user=visualops
loglevel=info
environment=PYTHON_EGG_CACHE=/visualops/.python-eggs,PYTHONPATH=/visualops/source/api/Source,INSTANT_HOME=/visualops
stdout_logfile=/visualops/log/supervisord/reqestworker8.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
redirect_stderr=true
directory=/visualops/source/api/
